# app.service
[Unit]
Description = app

[Service]
RemainAfterExit=true
ExecStop=/home/admins/app/start.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
