# app.service
[Unit]
Description = app

[Service]
Type=sipmle
ExecStart=/home/admins/app/LexusPokus

[Install]
WantedBy=multi-user.target
