# c#app.service
[Unit]
Description = c-app

[Service]
RemainAfterExit=true
ExecStop=/home/admins/app/LexusPokus
Type=oneshot

[Install]
WantedBy=multi-user.target
